# homeassistant-addon/ups2mqtt/rootfs/etc/services.d/ups2mqtt/run
#!/usr/bin/env sh
set -eu

export UPS2MQTT_OPTIONS_PATH="/data/ups2mqtt_internal_options.json"
export UPS2MQTT_LOG_LEVEL="$(python3 - <<'PY'
import json
from pathlib import Path

options = json.loads(Path('/data/options.json').read_text(encoding='utf-8'))
print(str(options.get('log_level', 'INFO')))
PY
)"

# Keep the Supervisor ingress prefix available for HA panel requests. The web
# handler detects direct-port requests and serves them from /.
if [ -n "${SUPERVISOR_INGRESS:-}" ]; then
  export UPS2MQTT_WEB_BASE_PATH="${SUPERVISOR_INGRESS}"
else
  export UPS2MQTT_WEB_BASE_PATH=""
fi

cd /usr/src/app
exec uv run python -m ups2mqtt.main
